# infra/sfw-sasuke.service
[Unit]
Description=SFW sasuke main app.

[Service]
ExecStart=/home/opc/sfw-sasuke/sfw-sasuke
Restart=always

[Install]
WantedBy=multi-user.target